# crontab
* * * * * /opt/cl0ck/updatedisplay.sh 2 >> /var/log/cl0ck.err
